# supervisord.conf
[supervisord]
logfile=/tmp/supervisord_privoxy.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord_privoxy.pid
nodaemon=true
silent=false
minfds=1024
minprocs=200
user=privoxy

[program:privoxy]
command=privoxy --no-daemon
directory=%(here)s
autostart=true
autorestart=true
redirect_stderr=True